# systemd/clownfisch.service
[Unit]
Description=🐠 Clownfischserver Telegram Bot
Documentation=https://github.com/YOUR_USERNAME/clownfischserver
After=network-online.target ollama.service
Wants=network-online.target
Requires=ollama.service

[Service]
Type=simple
User=clownfish
Group=clownfish
WorkingDirectory=/opt/clownfischserver
EnvironmentFile=/opt/clownfischserver/config/.env
ExecStart=/opt/clownfischserver/venv/bin/python3 /opt/clownfischserver/bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clownfisch

# Security
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
